# server/openclaw/systemd/todo-ilink-poller.timer
[Unit]
Description=Run Todo WeChat iLink message poll every 15 seconds

[Timer]
OnBootSec=20s
OnUnitActiveSec=15s
AccuracySec=5s
Persistent=true
Unit=todo-ilink-poller.service

[Install]
WantedBy=timers.target
