# server/openclaw/systemd/todo-openclaw-reminders.service
[Unit]
Description=Todo WeChat reminder scan
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/to_do
EnvironmentFile=/etc/todo-openclaw.env
ExecStart=/usr/bin/node /opt/to_do/server/openclaw/reminder-worker.mjs
